# riigikogu-radar/scripts/operatives/riigikogu-operatives.service
[Unit]
Description=Riigikogu Radar Operative Supervisor
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/riigikogu-radar/riigikogu-radar
ExecStart=/bin/bash -c 'cd /home/ubuntu/riigikogu-radar/riigikogu-radar && /usr/bin/npx tsx scripts/operatives/supervisor.ts'
Restart=always
RestartSec=60
StandardOutput=append:/home/ubuntu/riigikogu-radar/riigikogu-radar/logs/operatives/supervisor.log
StandardError=append:/home/ubuntu/riigikogu-radar/riigikogu-radar/logs/operatives/supervisor.log

# Environment
Environment=NODE_ENV=production
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
